# rootfs/etc/services.d/sshd/run
#!/bin/bash
# SSH server (optional, enabled via SSH_ENABLE=true)
source /etc/s6-overlay/lib/env-utils.sh
source_env_prefix SSH_

if [ "${SSH_ENABLE:-false}" != "true" ]; then
  echo "[sshd] SSH disabled (set SSH_ENABLE=true to enable)"
  s6-svc -Od .
  exit 0
fi

echo "[sshd] Starting SSH server..."

# Generate host keys if they don't exist
ssh-keygen -A 2>/dev/null || true

# Run sshd in foreground mode
exec /usr/sbin/sshd -D -e
