# rootfs/etc/services.d/openclaw/run
#!/bin/bash
# OpenClaw gateway service
source /etc/s6-overlay/lib/env-utils.sh
source_env_prefix OPENCLAW_

echo "[openclaw] Starting openclaw gateway..."
# Allow openclaw to use tailscale
if [ "${TAILSCALE_ENABLE:-false}" = "true" ]; then
    tailscale set --operator=openclaw
    echo "[openclaw] Starting with --tailscale serve"
    with_env_prefix --user openclaw OPENCLAW_ -- openclaw gateway --allow-unconfigured --tailscale serve
else
    with_env_prefix --user openclaw OPENCLAW_ -- openclaw gateway --allow-unconfigured
fi
