# rootfs/etc/services.d/moltbot/run
#!/command/with-contenv bash
# OpenClaw gateway service

echo "[moltbot] Starting moltbot gateway..."

# Run as moltbot user with proper environment
export HOME=/home/moltbot
exec s6-setuidgid moltbot bash -c '
  export HOME=/home/moltbot
  source /home/moltbot/.bashrc
  source /home/moltbot/.profile
  moltbot gateway --allow-unconfigured
' > /proc/1/fd/1 2>&1
