# rootfs/etc/services.d/crond/run
#!/bin/bash
# Start cron daemon in foreground
exec /usr/sbin/cron -f
